# files/systemd/both/tealdeer/tealdeer.service
[Unit]
Description=updates tldr local cache
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=tldr --update
