# files/systemd/both/tealdeer/tealdeer.timer
[Unit]
Description=updates tldr local cache

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1hours

[Install]
WantedBy=default.target
