# files/systemd/user/nnnplugins/nnnplugins.timer
[Unit]
Description=installs or updates nnn plugins

[Timer]
OnCalendar=weekly
Persistent=true
AccuracySec=1hours

[Install]
WantedBy=default.target
